# deploy/gruvboard.service
[Unit]
Description=GruvBoard - Homelab Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=gruvboard
Group=gruvboard
WorkingDirectory=/opt/gruvboard
ExecStart=/opt/gruvboard/.venv/bin/uvicorn backend.main:app --host 127.0.0.1 --port 9500
Restart=on-failure
RestartSec=5
Environment=GRUVBOARD_CONFIG=/opt/gruvboard/config.toml

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/gruvboard/data
PrivateTmp=true

[Install]
WantedBy=multi-user.target
